# build_dirs/001/crond.init
#!/sbin/openrc-run

name="busybox crond"
command="/usr/sbin/crond"
command_args="${CRON_OPTS}"
pidfile="/var/run/crond.pid"

depend() {
	need localmount
	use logger
}

depend() {
	provide logger
}
